# deploy/systemd/cabo-backup.timer
[Unit]
Description=Run the Cabo database backup every day

[Timer]
OnCalendar=*-*-* 03:15:00 UTC
RandomizedDelaySec=15m
Persistent=true
Unit=cabo-backup.service

[Install]
WantedBy=timers.target
